# conf/etc/systemd/system/igmpproxy.service.d/override.conf
[Unit]
PartOf=networking.service
After=networking.service
